# kafka/supervisor/kafka.conf
[program:kafka]
command=/usr/bin/start-kafka.sh
autostart=true
autorestart=true
environment=ADVERTISED_HOST="%(ADVERTISED_HOST)s"